# deploy/analizador-tdr.service
[Unit]
Description=Vigilante SEACE - Analizador TDR (FastAPI)
After=network.target
# Límites de reinicio — DEBE estar en [Unit] para systemd <248
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/vigilante-seace/analizador-tdr

ExecStart=/var/www/vigilante-seace/analizador-tdr/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8001

# Uvicorn responde bien a SIGTERM — graceful shutdown
KillSignal=SIGTERM
TimeoutStopSec=15

Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/vigilante-seace/analizador-tdr.log
StandardError=append:/var/log/vigilante-seace/analizador-tdr-error.log

# Seguridad
NoNewPrivileges=true
ReadWritePaths=/var/www/vigilante-seace/analizador-tdr /var/log/vigilante-seace /tmp

[Install]
WantedBy=multi-user.target
