# deploy/telegram-bot.service
[Unit]
Description=Vigilante SEACE - Telegram Bot Listener
After=network.target mysql.service
Wants=mysql.service
# Límites de reinicio (máx 5 en 5 min) — DEBE estar en [Unit] para systemd <248
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/vigilante-seace

ExecStart=/usr/local/php82/bin/php artisan telegram:listen

# Parada limpia: SIGTERM al proceso principal, SIGKILL al cgroup tras timeout
# El long-polling HTTP bloquea hasta 30s — dejamos 35s antes del SIGKILL
# NOTA: NO usar --isolated aquí; systemd ya garantiza instancia única
#       y --isolated deja locks stale si el proceso muere con SIGKILL
KillSignal=SIGTERM
KillMode=mixed
TimeoutStopSec=35

Restart=on-failure
RestartSec=30
StandardOutput=append:/var/log/vigilante-seace/telegram-bot.log
StandardError=append:/var/log/vigilante-seace/telegram-bot-error.log

# Seguridad
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/www/vigilante-seace/storage /var/log/vigilante-seace /tmp

[Install]
WantedBy=multi-user.target
