# deploy/telegram-admin-bot.service
[Unit]
Description=Vigilante SEACE - Telegram Admin Bot Listener
After=network.target mysql.service
Wants=mysql.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/vigilante-seace

ExecStart=/usr/local/php82/bin/php artisan telegram:admin-listen

KillSignal=SIGTERM
KillMode=mixed
TimeoutStopSec=20

Restart=on-failure
RestartSec=30
StandardOutput=append:/var/log/vigilante-seace/telegram-admin-bot.log
StandardError=append:/var/log/vigilante-seace/telegram-admin-bot-error.log

# Seguridad
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/www/vigilante-seace/storage /var/log/vigilante-seace /tmp

[Install]
WantedBy=multi-user.target
