# deploy/qa/analizador-tdr-qa.service
[Unit]
Description=Vigilante SEACE QA - Analizador TDR (FastAPI)
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/vigilante-seace-qa/analizador-tdr

ExecStart=/var/www/vigilante-seace-qa/analizador-tdr/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8002

KillSignal=SIGTERM
TimeoutStopSec=15

Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/vigilante-seace-qa/analizador-tdr.log
StandardError=append:/var/log/vigilante-seace-qa/analizador-tdr-error.log

NoNewPrivileges=true
ReadWritePaths=/var/www/vigilante-seace-qa/analizador-tdr /var/log/vigilante-seace-qa /tmp

[Install]
WantedBy=multi-user.target
